# crontab
0 18 * * * /app/cleanup.sh >> /var/log/cron.log 2>&1
